# overlay/airootfs/etc/udev/rules.d/98-monitor-plugin.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{TRIGGERED_BY_UDEV}="1", RUN+="/etc/p3ng0s/display.sh"
