# overlay/airootfs/etc/systemd/system/p3ng0s-certs-watcher.path
[Unit]
Description=Watch for new certificates in loot certs folder
After=p3ng0s-mount-loot.service
Requires=p3ng0s-mount-loot.service

[Path]
PathExistsGlob=/home/p4p1-live/loot/certs/*.cer
PathChanged=/home/p4p1-live/loot/certs
Unit=p3ng0s-certs-worker.service

[Install]
WantedBy=multi-user.target
