# overlay/airootfs/etc/systemd/system/p3ng0s-dropbox-nethide.service
[Unit]
Description=MAC Address Rotation

[Service]
Type=oneshot
ExecStart=/etc/p3ng0s/dropbox/nethide.sh
User=root
